# detections/splunk-imports/_quarantine/linux-auditd-possible-append-cronjob-entry-on-existing-cronjob-file.yaml
id: splunk-security-content-fea71cf0-fa10-4ef6-9202-9682b2e0c477
name: Linux Auditd Possible Append Cronjob Entry On Existing Cronjob File
description: 'The following analytic detects potential tampering with cronjob files
  on a Linux system.

  It leverages logs from Linux Auditd, focusing on events of type PATH or CWD.

  This activity could be significant because adversaries often use it for persistence
  or privilege escalation.

  Correlate this with related EXECVE or PROCTITLE events to identify the process or
  user responsible for the access or modification.

  If confirmed malicious, this could allow attackers to execute unauthorized code
  automatically, leading to system compromises and unauthorized data access, thereby
  impacting business operations and data integrity.'
severity: medium
enabled: false
tags:
  mitre:
  - T1053.003
  categories:
  - endpoint
references:
- https://attack.mitre.org/techniques/T1053/003/
- https://blog.aquasec.com/threat-alert-kinsing-malware-container-vulnerability
- https://www.intezer.com/blog/research/kaiji-new-chinese-linux-malware-turning-to-golang/
logsource:
  product: splunk
detection:
  splunk_spl: "`linux_auditd` (type=PATH OR type=CWD)\n| rex \"msg=audit\\([^)]*:(?<audit_id>\\\
    d+)\\)\"\n\n| stats\n  values(type) as types\n  values(name) as names\n  values(nametype)\
    \ as nametype\n  values(cwd) as cwd_list\n  values(_time) as event_times\n  by\
    \ audit_id, host\n\n| eval current_working_directory = coalesce(mvindex(cwd_list,\
    \ 0), \"N/A\")\n| eval candidate_paths = mvmap(names, if(match(names, \"^/\"),\
    \ names, current_working_directory + \"/\" + names))\n| eval matched_paths = mvfilter(match(candidate_paths,\
    \ \"/etc/cron.*|.*/cron/.*|/etc/anacrontab.*\"))\n| eval match_count = mvcount(matched_paths)\n\
    | eval reconstructed_path = mvindex(matched_paths, 0)\n| eval e_time = mvindex(event_times,\
    \ 0)\n| where match_count > 0\n| rename host as dest\n\n| stats count min(e_time)\
    \ as firstTime max(e_time) as lastTime\n  values(nametype) as nametype\n  by current_working_directory\n\
    \     reconstructed_path\n     match_count\n     dest\n     audit_id\n\n| `security_content_ctime(firstTime)`\n\
    | `security_content_ctime(lastTime)`\n| `linux_auditd_possible_append_cronjob_entry_on_existing_cronjob_file_filter`"
provenance:
  source: splunk/security_content
  source_id: fea71cf0-fa10-4ef6-9202-9682b2e0c477
  source_commit: 4d4c7ee
  license: Apache-2.0
  license_url: https://github.com/splunk/security_content/blob/develop/LICENSE
  imported_at: '2026-05-04'
  imported_by: splunk_importer
  upstream_path: detections/endpoint/linux_auditd_possible_append_cronjob_entry_on_existing_cronjob_file.yml
notes:
  splunk_status: production
  splunk_data_source: []
  quarantine_reason: raw SPL — needs translation to AiSOC detection schema
