# detections/car-imports/endpoint/car-2013-08-001.yaml
id: mitre-car-car-2013-08-001
name: Execution with schtasks
description: The Windows built-in tool `schtasks.exe` provides the creation, modification,
  and running of [scheduled tasks](https://attack.mitre.org/techniques/T1053) on a
  local or remote computer. It is provided as a more flexible alternative to `at.exe`,
  described in [CAR-2013-05-004](../CAR-2013-05-004). Although used by adversaries,
  the tool is also legitimately used by administrators, scripts, and software configurations.
  The scheduled tasks tool can be used to gain [Persistence](https://attack.mitre.org/tactics/TA0003)
  and can be used in combination with a [Lateral Movement](https://attack.mitre.org/tactics/TA0008)
  technique to remotely gain [execution](https://attack.mitre.org/tactics/TA0002).
  Additionally, the command has parameters to specify the user and password responsible
  for creating the task, as well as the user and password combination that the task
  will run as. The `/s` flag specifies the remote system on which the task should
  be scheduled, usually indicating [Lateral Movement](https://attack.mitre.org/tactics/TA0008).
severity: medium
enabled: false
tags:
  mitre:
  - T1053
  categories:
  - endpoint
references: []
logsource:
  category: Process
detection:
  implementations:
  - name: ''
    type: pseudocode
    code: 'process = search Process:Create

      schtasks = filter process where (exe == "schtasks.exe")

      output schtasks

      '
  - name: ''
    type: DNIF
    code: '_fetch * from event where $LogName=WINDOWS-SYSMON AND $EventID=1 AND $App=schtasks.exe
      AND $Process=regex(.*(\/create|\/run|\/query|\/delete|\/change|\/end).*)i limit
      100

      '
  - name: ''
    type: LogPoint
    code: 'norm_id=WindowsSysmon event_id=1 image="*\schtasks.exe" command IN ["*/create*",
      "*/run*", "*/query*", "*/delete*", "*/change*", "*/end*"]

      '
provenance:
  source: mitre-attack/car
  source_id: CAR-2013-08-001
  source_commit: 1b922fe
  license: Apache-2.0
  license_url: https://github.com/mitre-attack/car/blob/master/LICENSE.txt
  imported_at: '2026-05-04'
  imported_by: car_importer
  upstream_path: analytics/CAR-2013-08-001.yaml
notes:
  car_status: needs_implementation
