# detections/chronicle-imports/_quarantine/suspicious-scheduled-task.yaml
id: chronicle-detection-rules-suspicious-scheduled-task
name: suspicious_scheduled_task
description: suspicious_scheduled_task
severity: medium
enabled: false
tags:
  mitre: []
  categories:
  - endpoint
references: []
logsource:
  product: chronicle
detection:
  chronicle_yaral: "rule suspicious_scheduled_task {\n meta:\n    author = \"Emir\
    \ Erdogan\"\n    description = \"Detection of suspicious scheduled tasks  License:\
    \ https://github.com/Neo23x0/sigma/blob/master/LICENSE.Detection.Rules.md.\"\n\
    \    reference = \"https://tdm.socprime.com/tdm/info/1ULlEDq8oA0b\"\n    version\
    \ = \"0.01\"\n    created = \"2021-03-09\"\n    product = \"windows\"\n    service\
    \ = \"sysmon\"\n    mitre = \"T1053, T1168\"\n\n  events:\n(re.regex($selection.target.process.command_line,\
    \ `.*schtasks\\.exe /create /sc MINUTE /tn SystemSoundsServices /tr \\\"regsvr32\\\
    .exe\\\".*`) or re.regex($selection.target.process.command_line, `.*/s /n /u /i:http:.*`))\n\
    \n  condition:\n    $selection\n}"
  chronicle_rule_name: suspicious_scheduled_task
provenance:
  source: chronicle/detection-rules
  source_id: suspicious_scheduled_task
  source_commit: 74dd490
  license: Apache-2.0
  license_url: https://github.com/chronicle/detection-rules/blob/main/LICENSE
  imported_at: '2026-05-04'
  imported_by: chronicle_importer
  upstream_path: rules/_deprecated/soc_prime_rules/threat_hunting/sysmon/suspicious_scheduled_task.yaral
notes:
  chronicle_meta: {}
  quarantine_reason: raw YARA-L — needs translation to AiSOC schema
