# detections/chronicle-imports/_quarantine/suspicious-schtasks-creation-possible-windows-0day-lpe-aka-polarbear-by-sandboxe.yaml
id: chronicle-detection-rules-suspicious-schtasks-creation-possible-windows-0day-lpe-aka-polarbear-by-sandboxe
name: suspicious_schtasks_creation_possible_windows_0day_lpe_aka_polarbear_by_sandboxescaper
description: suspicious_schtasks_creation_possible_windows_0day_lpe_aka_polarbear_by_sandboxescaper
severity: medium
enabled: false
tags:
  mitre: []
  categories:
  - endpoint
references: []
logsource:
  product: chronicle
detection:
  chronicle_yaral: "rule suspicious_schtasks_creation_possible_windows_0day_lpe_aka_polarbear_by_sandboxescaper\
    \ {\n meta:\n    author = \"Roman Ranskyi\"\n    description = \"Possible task\
    \ scheduler .job import arbitrary DACL write.  License: https://github.com/Neo23x0/sigma/blob/master/LICENSE.Detection.Rules.md.\"\
    \n    reference = \"https://tdm.socprime.com/tdm/info/Q9PtOg89jvkq\"\n    version\
    \ = \"0.01\"\n    created = \"2019/05/22\"\n    product = \"windows\"\n    service\
    \ = \"process_creation\"\n    mitre = \"Execution, Privilege Escalation, t1053\"\
    \n\n  events:\n(($selection.metadata.product_event_type = \"1\" or $selection.metadata.product_event_type\
    \ = \"4688\") and $selection.target.process.command_line = \"schtasks\" and $selection.target.process.command_line\
    \ = \"/change\" and $selection.target.process.command_line = \"/TN\" and $selection.target.process.command_line\
    \ = \"/RU\" and $selection.target.process.command_line = \"/RP\")\n\n  condition:\n\
    \    $selection\n}"
  chronicle_rule_name: suspicious_schtasks_creation_possible_windows_0day_lpe_aka_polarbear_by_sandboxescaper
provenance:
  source: chronicle/detection-rules
  source_id: suspicious_schtasks_creation_possible_windows_0day_lpe_aka_polarbear_by_sandboxescaper
  source_commit: 74dd490
  license: Apache-2.0
  license_url: https://github.com/chronicle/detection-rules/blob/main/LICENSE
  imported_at: '2026-05-04'
  imported_by: chronicle_importer
  upstream_path: rules/_deprecated/soc_prime_rules/threat_hunting/microsoft_sysmon/suspicious_schtasks_creation__possible_windows_0day_lpe_aka__polarbear__by_sandboxescaper.yaral
notes:
  chronicle_meta: {}
  quarantine_reason: raw YARA-L — needs translation to AiSOC schema
