# detections/endpoint/linux-anacron-job-add.yaml
id: det-endpoint-134
name: Anacron Job File Added
description: AiSOC v1 curated detection. Triggers on the endpoint signal described by 'Anacron Job File
  Added'. Watch the 1 documented false-positive case before tuning.
version: 1.0.0
severity: medium
tags:
- mitre.attack.t1053.003
- tlp.white
category: endpoint
log_source:
  product: linux
  service: auditd
detection:
  condition: |
    syscall == "openat"
    AND path CONTAINS_ANY ["/etc/anacrontab", "/etc/cron.daily/", "/etc/cron.weekly/"]
    AND exe NOT ENDSWITH_ANY ["/dpkg", "/rpm"]
false_positives:
- Sanctioned automation/IaC pipeline
playbook: tpl-persistence
enabled: true
author: AiSOC
created: '2026-05-03'
modified: '2026-05-03'
