# detections/endpoint/linux-cron-d-write.yaml
id: det-endpoint-133
name: File Written Into /etc/cron.d Outside Package Mgr
description: AiSOC v1 curated detection. Triggers on the endpoint signal described by 'File Written Into
  /etc/cron.d Outside Package Mgr'. Watch the 1 documented false-positive case before tuning.
version: 1.0.0
severity: high
tags:
- mitre.attack.t1053.003
- tlp.white
category: endpoint
log_source:
  product: linux
  service: auditd
detection:
  condition: |
    syscall == "openat"
    AND path CONTAINS "/etc/cron.d/"
    AND exe NOT ENDSWITH_ANY ["/dpkg", "/rpm", "/apt-get", "/yum", "/dnf"]
false_positives:
- Sanctioned automation/IaC pipeline
playbook: tpl-persistence
enabled: true
author: AiSOC
created: '2026-05-03'
modified: '2026-05-03'
