# detections/endpoint/linux-systemd-service-create-untrusted.yaml
id: det-endpoint-131
name: Systemd Unit File Written By Non-Package Process
description: AiSOC v1 curated detection. Triggers on the endpoint signal described by 'Systemd Unit File
  Written By Non-Package Process'. Watch the 2 documented false-positive cases before tuning.
version: 1.0.0
severity: high
tags:
- mitre.attack.t1543.002
- tlp.white
category: endpoint
log_source:
  product: linux
  service: auditd
detection:
  condition: |
    syscall == "openat"
    AND path CONTAINS "/etc/systemd/system/"
    AND exe NOT ENDSWITH_ANY ["/dpkg", "/rpm", "/yum", "/apt-get", "/zypper"]
false_positives:
- Sanctioned automation/IaC pipeline
- Configuration management agent
playbook: tpl-persistence
enabled: true
author: AiSOC
created: '2026-05-03'
modified: '2026-05-03'
