# detections/endpoint/win-task-create-from-untrusted-process.yaml
id: det-endpoint-044
name: Scheduled Task Created From Office or Browser
description: AiSOC v1 curated detection. Triggers on the endpoint signal described by 'Scheduled Task
  Created From Office or Browser'. Watch the 1 documented false-positive case before tuning.
version: 1.0.0
severity: high
tags:
- mitre.attack.t1053.005
- tlp.white
category: endpoint
log_source:
  product: windows
  service: task-scheduler
detection:
  condition: |
    event_id == 4698
    AND creator_process ENDSWITH_ANY ["winword.exe", "excel.exe", "outlook.exe", "chrome.exe", "msedge.exe"]
false_positives:
- Sanctioned automation/IaC pipeline
playbook: tpl-persistence
enabled: true
author: AiSOC
created: '2026-05-03'
modified: '2026-05-03'
