# detections/endpoint/win-task-runs-encoded-powershell.yaml
id: det-endpoint-045
name: Scheduled Task Runs Encoded PowerShell
description: AiSOC v1 curated detection. Triggers on the endpoint signal described by 'Scheduled Task
  Runs Encoded PowerShell'. Watch the 1 documented false-positive case before tuning.
version: 1.0.0
severity: high
tags:
- mitre.attack.t1053.005
- mitre.attack.t1027
- tlp.white
category: endpoint
log_source:
  product: windows
  service: task-scheduler
detection:
  condition: |
    event_id == 4698
    AND task_command CONTAINS "-encodedcommand"
false_positives:
- Vendor agent that uses encoded payloads (allow-listed)
playbook: tpl-execution
enabled: true
author: AiSOC
created: '2026-05-03'
modified: '2026-05-03'
